# deploy/torque_influx.conf
# simple uWSGI script

description "uwsgi tiny instance"
start on runlevel [2345]
stop on runlevel [06]

script

  exec influxd

  echo 'torque influx up'
end script

respawn
